# .config/systemd/user/keepassxc.service
[Unit]
Description=KeepassXC Password Manager

[Service]
Environment="QT_QPA_PLATFORM=wayland"
ExecStart=/usr/bin/keepassxc
